# .config/systemd/user/git.service
[Unit]
Description=auto pull and push git

[Service]
Type=forking
ExecStart=cd; git pull;
ExecStop=cd; git push;

[Install]
WantedBy=default.target

